# .config/systemd/user/unity7.service.d/lowgfx.conf
[Service]
Environment=UNITY_LOW_GFX_MODE=1
